# ansible/playbooks/templates/security-updates.service.j2
[Unit]
Description=Apply security updates using dnf
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/dnf -y --security upgrade
