# ansible/playbooks/templates/security-updates.timer.j2
[Unit]
Description=Schedule security updates (fallback)

[Timer]
OnCalendar={{ patch_window }}
Persistent=true

[Install]
WantedBy=timers.target
