# ansible/roles/patching/templates/rhel/security-updates.timer.j2
[Unit]
Description=Schedule security updates (fallback)

[Timer]
OnCalendar={{ patching_rhel_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target

