# ansible/roles/patching/templates/rhel/security-updates.service.j2
[Unit]
Description=Apply security updates using {{ _dnf_cmd }}

[Service]
Type=oneshot
ExecStart=/usr/bin/{{ _dnf_cmd }} -y --security upgrade

